# root/etc/supervisor.d/slurmd.conf
[program:slurmd]
command=sh -c "sleep 10; /opt/slurm/sbin/slurmd -D -vvvvvv"
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autostart=true
autorestart=false
startretries=0
stopwaitsecs=1
exitcodes=0,1,2
